# deploy/ec2/gunicorn.socket
[Unit]
Description=Gunicorn socket for Siyasi News

[Socket]
ListenStream=/run/siyasi_news.sock

[Install]
WantedBy=sockets.target
